# apps/kodi/src/bash-docker.service
[Unit]
Description=Bash Docker Service

[Service]
ExecStart=/usr/src/app/start.sh

[Install]
WantedBy=basic.target
